# deploy/monitor-aforo.service
# =============================================================================
#  Servicio systemd para el backend de Monitor de Aforo
#  Ubicación en el servidor: /etc/systemd/system/monitor-aforo.service
#
#  Activar:
#    sudo cp monitor-aforo.service /etc/systemd/system/
#    sudo systemctl daemon-reload
#    sudo systemctl enable monitor-aforo
#    sudo systemctl start  monitor-aforo
# =============================================================================

[Unit]
Description=Monitor de Aforo — FastAPI backend
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data

# Directorio raíz del proyecto en el servidor
WorkingDirectory=/opt/monitor-aforo/backend

# Carga las variables de entorno desde .env
EnvironmentFile=/opt/monitor-aforo/backend/.env

# Arranca uvicorn con el entorno virtual del proyecto
ExecStart=/opt/monitor-aforo/backend/venv/bin/uvicorn main:app \
          --host 127.0.0.1 \
          --port 8000 \
          --workers 1

# Reinicio automático si falla (con espera de 5 s entre intentos)
Restart=on-failure
RestartSec=5s

# Límites de recursos
LimitNOFILE=65536
TimeoutStopSec=30

# Logging: los logs van a journald (ver con: journalctl -u monitor-aforo -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=monitor-aforo

[Install]
WantedBy=multi-user.target
